# run-poly.service
[Unit]
Description=Crypto Trading Bot with NGINX Proxy
After=network.target nginx.service
Requires=nginx.service

[Service]
User=admin
Group=admin
WorkingDirectory=/home/admin/poly_22_headless
ExecStart=/usr/bin/bash /home/admin/poly_22_headless/run.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/admin/poly_22_headless/run.log
StandardError=append:/home/admin/poly_22_headless/run.log

# 环境变量
Environment=FLASK_HOST=127.0.0.1
Environment=FLASK_PORT=5000
Environment=FLASK_ENV=production

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/admin/poly_22_headless

# 资源限制
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target