# run-poly.service
[Unit]
Description=Crypto Trading Bot with NGINX Proxy
After=network.target nginx.service
Requires=nginx.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/poly_22_headless
Environment=PATH=%h/poly_22_headless/venv/bin:/usr/bin:/bin
Environment=FLASK_HOST=127.0.0.1
Environment=FLASK_PORT=5000
Environment=FLASK_ENV=production
ExecStart=%h/poly_22_headless/venv/bin/python %h/poly_22_headless/crypto_trader.py
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]  
WantedBy=multi-user.target